# configs/conf.d/filebeat-lumen.conf
[program:filebeat-lumen]
command=filebeat -e -c /filebeat.yml
autostart=true
autorestart=true
priority=1
stdout_events_enabled=true
stderr_events_enabled=true